# config/tenantfirstaid-backend.service
[Unit]
Description=Gunicorn + uv for tenantfirstaid Flask API
After=network.target

[Service]
User=root
Group=root
WorkingDirectory=/var/www/tenantfirstaid/backend

# ── main line ─────────────────────────────────────────────
ExecStart=/root/.local/bin/uv run --no-sync gunicorn -w 1 -b unix:/run/tenantfirstaid.sock tenantfirstaid.app:app

Restart=on-failure
KillSignal=SIGQUIT
TimeoutStopSec=10

[Install]
WantedBy=multi-user.target
